# templates/profile.d/slurm.sh.erb
#!/bin/sh
#
# File managed by Puppet!
#

<%- if scope['slurm::profiled_add_path'] -%>
export PATH="<%= scope['slurm::source_install_prefix'] %>/bin:$PATH"
<%- end -%>
export SLURM_CONF="<%= scope['slurm::slurm_conf_path'] %>"
