# templates/profile.d/slurm.csh.erb
#!/bin/csh
#
# File managed by Puppet!
#

<%- if scope['slurm::profiled_add_path'] -%>
setenv PATH "<%= scope['slurm::source_install_prefix'] %>/bin:$PATH"
<%- end -%>
<%- if !scope['slurm::configless'] -%>
setenv SLURM_CONF "<%= scope['slurm::slurm_conf_path'] %>"
<%- end -%>
